# rpi-amp-standby.service
[Install]
WantedBy=multi-user.target
After=pulseaudio.target

[Service]
ExecStart=/usr/local/bin/rpi-amp-standby
